# systemd/robot-web-launcher.service
[Unit]
Description=Robot Web Launcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nvidia1
WorkingDirectory=/data/sinuo_project/robot_web_launcher
Environment=PYTHONUNBUFFERED=1
Environment=ROBOT_LAUNCHER_CONFIG=/data/sinuo_project/robot_web_launcher/config/modules.yaml
Environment=ROBOT_LAUNCHER_AUTOSTART=0
Environment=ROBOT_LAUNCHER_PYTHON=/data/sinuo_project/robot_web_launcher/.venv/bin/python
ExecStart=/usr/bin/env bash scripts/run_dev.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
